# infra/systemd/youtube-upload-worker.service
[Unit]
Description=Youtube Upload Lush Worker Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/youtube-upload-lush
EnvironmentFile=/etc/youtube-upload-worker.env
Environment=HOME=/opt/youtube-upload-lush-runtime/system-home
ExecStart=/opt/youtube-upload-lush/.venv/bin/python -m workers.agent.main
Restart=always
RestartSec=5
User=ytworker
Group=ytworker

[Install]
WantedBy=multi-user.target
